# supervisor-tsdb.conf
[program:opentsdb]
command=/opt/start/start_opentsdb.sh
stopwaitsecs=60
