# supervisor-hbase.conf
[program:hbase]
command=/opt/start/start_hbase.sh
stopwaitsecs=60
